# shimming/alienkeyboardservice/alienkeyboard-maliit-shim.service
[Unit]
Description=Aliendalvik keyboard DBus shim provider
After=dbus.socket graphical-session.target
Requires=dbus.socket graphical-session.target

[Service]
Type=simple
Restart=always
RestartSec=5
ExecStart=/usr/bin/alienkeyboard-maliit-shim/alienkeyboard-maliit-shim
